# files/4-rinkhals/etc/profile.d/rinkhals.sh
export HOME=/useremain/home/rinkhals
export PATH=$PATH:$HOME:$HOME/bin:$HOME/usr/bin:$HOME/sbin:$HOME/usr/sbin
export LD_LIBRARY_PATH=$HOME/usr/lib:$HOME/lib:$LD_LIBRARY_PATH

export PATH=/usr/libexec/gcc/arm-buildroot-linux-uclibcgnueabihf/11.4.0:$PATH
export CC=/usr/bin/gcc
